# resources/lib/systemd/system/intel-rapl-perms.service
[Unit]
Description=Set rapl permissions for games

[Service]
Type=oneshot
ExecStart=/usr/bin/intel-rapl-perms

[Install]
WantedBy=graphical.target
